# systemd/user/dropbox.service.d/override.conf
[Service]
Environment=DISPLAY=:0
Type=forking
PIDFile=/home/lgunsch/.dropbox/dropbox.pid
